# autopush.cron
15,30,45,59 * * * * cd /home/pi/Documents/python/vc-info/&&python3 /home/pi/Documents/python/vc-info/VancraftInfo.py
15,30,45,59 * * * * cd /home/pi/Documents/python/vc-info/&&/home/pi/Documents/python/vc-info/push.sh
