# autopull.cron
15,30,45,59 * * * * cd /home/pi/vc-info/&&python3 /home/pi/vc-info/VancraftInfo.py
15,30,45,59 * * * * cd /home/pi/vc-info/&&/home/pi/vc-info/pull.sh
